# com.saber.goultra.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.saber.goultra</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO_DIR__/trigger.sh</string>
    </array>

    <key>StartOnMount</key>
    <true/>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/launchd-stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
